# ports/net-mgmt/bird_exporter/files/bird_exporter.in
#!/bin/sh
#
# PROVIDE: bird_exporter
# REQUIRE: LOGIN bird
# KEYWORD: shutdown

. /etc/rc.subr

name="bird_exporter"
rcvar="bird_exporter_enable"
desc="Prometheus exporter for BIRD routing daemon"

load_rc_config ${name}

: ${bird_exporter_enable:="NO"}
: ${bird_exporter_user:="nobody"}
: ${bird_exporter_socket:="/var/run/bird/bird.ctl"}
: ${bird_exporter_args:="-bird.v2 -bird.socket ${bird_exporter_socket}"}
: ${bird_exporter_listen:="9324"}

pidfile="/var/run/${name}.pid"
command="/usr/sbin/daemon"
command_args="-P ${pidfile} -r -f /usr/local/bin/${name} ${bird_exporter_args} -web.listen-address :${bird_exporter_listen}"

run_rc_command "$1"
